# etc/systemd/system/ssg-services.target
[Unit]
Description=services for static site
Requires=multi-user.target
After=multi-user.target
AllowIsolate=yes
